# services/oci-cn-auth.service.in
#
# oci-cn-auth. service providing configuration of wpa_supplicant for specified interface
# to authenticate to OCI cluster network.
#
[Unit]
Description=OCI Cluster Network Authentication: renew WPA certificates (__VERSION__)
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/mkdir -p /var/run/oci-cn-auth
ExecStartPre=/bin/chmod 0755 /var/run/oci-cn-auth
ExecStart=/opt/oci-hpc/oci-cn-auth/oci_cn_auth --start
ExecStop=/opt/oci-hpc/oci-cn-auth/oci_cn_auth --stop

[Install]
WantedBy=multi-user.target
